# conf/hw_tuning/hft.slice
[Unit]
Description=HFT hot-path reserved CPUs
Before=slices.target

[Slice]
AllowedCPUs=4-7 20-23
