# conf/feeder.service
[Unit]
Description=HFT Market Data Feeder
After=network.target

[Service]
Type=simple
User=hft-user
Group=hft-group
WorkingDirectory=/opt/hft
ExecStart=/opt/hft/bin/feeder/feeder
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hft-feeder

[Install]
WantedBy=multi-user.target
